# systemd/changelog-monitor.service
[Unit]
Description=Changelog Monitor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/gitdata/ollama_scanner
ExecStart=/bin/bash /root/gitdata/ollama_scanner/scripts/changelog_monitor.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target